# sleep.cron
23 14 * * * env DISPLAY=:0 java -jar /home/alecsnyder/Documents/Java/shutdown.jar 
