# sleep.cron
19 1 * * 1 wall < /home/alecsnyder/Documents/git/cs162/sleepTime/warning.txt
