# sleep.cron
32 1 * * 3 wall < /home/alecsnyder/Documents/git/cs162/sleepTime/warning.txt
